# deploy/renovation-api.service.d/10-env.conf
# systemd drop-in for the RenoHub API — environment only, no host-specific paths.
# Safe to `sudo cp` onto any host without clobbering the base unit's User/paths.
#   sudo mkdir -p /etc/systemd/system/renovation-api.service.d
#   sudo cp deploy/renovation-api.service.d/10-env.conf /etc/systemd/system/renovation-api.service.d/
#   sudo systemctl daemon-reload && sudo systemctl restart renovation-api
#
# RENOHUB_ENV=prod arms the fail-fast config guard: the app refuses to boot if
# JWT_SECRET is the dev default / under 32 chars or COOKIE_SECURE is false.
[Service]
Environment=RENOHUB_ENV=prod
